# package/base-files/files/etc/uci-defaults/14_vyacht_terrible_workaround_luci_index
#!/bin/sh

cp /www/index-vyacht.html /www/index.html 

exit 0
